# tg-private-bot.service
[Unit]
Description=Falcon Telegram Bot
After=network.target

[Service]
User=%i
WorkingDirectory=/home/%i/falcon-tg
ExecStart=/home/%i/falcon-tg/venv/bin/python /home/%i/falcon-tg/bot.py
EnvironmentFile=/home/%i/falcon-tg/.env
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
